# conf/systemd.service
[Unit]
Description=CryptPad is an encrypted PAD Creator
After=syslog.target
After=network.target

[Service]
Type=simple
User=__APP__
Group=__APP__
ExecStartPre=__NODEJS__
WorkingDirectory=__FINALPATH__
Environment="PATH=__ENV_PATH__"
ExecStart=/opt/node_n/bin/node server.js
Restart=always

[Install]
WantedBy=multi-user.target
