# conf/systemd.service
[Unit]
Description=CryptPad is an encrypted PAD Creator
After=syslog.target
After=network.target

[Service]
Type=simple
User=__APP__
Group=__APP__
ExecStartPre=__NODEJS__
WorkingDirectory=__FINALPATH__
ExecStart=__NODE__ server.js
Restart=always
Environment=USER=__APP__ HOME=__FINALPATH__

[Install]
WantedBy=multi-user.target
